# analyst/dist/templates/_base/config/systemd/analyst.service
###  THIS FILE IS A PLACEHOLDER AND A TEST
[Unit]
Description=Job that runs Zorroa Analyst in daemon mode

[Service]
User=computeruser
ExecStart=/home/computeruser/analyst/bin -d --pidfile=/home/computeruser/analyst/logs/analyst.log
Restart=on-abort


[Install]
WantedBy=multi-user.target
